# clodforest.service
[Unit]
Description=ClodForest Coordinator
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/ClodForest/lc_src
ExecStart=/bin/bash -c 'source venv/bin/activate && python clodforest.py'
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=REPO_PATH=./state
Environment=ENABLE_OAUTH2=true

[Install]
WantedBy=multi-user.target
